# sys_files/deku_daemon.service
[Unit]
Description=Deku Daemon can be used to route SMS messages via multiple SMS capable devices

